# cfengine/files/etc/cron.daily/upd-local-nss-db
#!/bin/sh
# Updated cached credentials database 
/usr/sbin/nss_updatedb ldap
