# cfengine/files/etc/cron.daily/upd_aliases_db
#!/bin/sh
# Update and make sure /etc/aliases.db exists
/usr/bin/newaliases