# cfengine/files/etc/cron.d/backup_cyso_2_bit
# rsync mrxbak0 to mexbak0
12 9 * * *      root /usr/local/sbin/backup_cyso_2_bit.sh
