# cfengine/files/etc/cron.d/m4n-mrhpgdb5
#
# This file is managed by cfengine
#
MAILTO=application@m4n.nl

# dump and restore production databases on a separate standby server 
0 8,9,10,11,12,13,14,15,16,17,18 * * * root time /m4n/bin/nightly_cron_queries.sh bench

# EOF