# cfengine/files/etc/cron.d/m4n-sat
#
# This file is managed by cfengine
#
MAILTO=application@m4n.nl

# backup tomcat log files
00 3    * * *   root    /m4n/bin/backup_tomcat.sh

# EOF