# cfengine/files/etc/cron.d/moxsts-logstats
#! /bin/bash

0 0   * * *     root    cd /data/logstats/;./sync-logs.sh;./generate_webalizer.sh;./generate_stats.sh

# EOF