# cfengine/files/etc/cron.d/ntp-check
#
# This file is managed by cfengine
#

# Run every half hour.
30 * * * * 	root	/usr/local/sbin/check_ntpd_xenU.sh

# EOF