# cfengine/files/etc/cron.d/rbackup-office
#
# This file is managed by cfengine
#
# Backups this server
4 32		* * *	root	/etc/rbackup/rbackup.sh s


# EOF